# files/etc/uci-defaults/99-vpn-enable
#!/bin/sh
/etc/init.d/vpn enable
rm -f /etc/uci-defaults/99-vpn-enable
